# .config/systemd/user/alc.target
[Unit]
Description=Starts stuff on login
Requires=default.target
Conflicts=rescue.service rescue.target
After=default.target
AllowIsolate=yes
