# .config/sytemd/user/mpd.service.d/override.conf
[Service]
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectControlGroups=true
ReadWritePaths=/home/alc/.config/mpd
RemoveIPC=true
RuntimeDirectory=mpd