# .config/sytemd/user/redshift-gtk.service.d/override.conf
[Service]
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
