# .config/sytemd/user/keybase.service.d/override.conf
[Service]
RuntimeDirectory=keybase
ProtectKernelTunables=true
RestrictNamespaces=true
ProtectSystem=strict
ReadWritePaths=/opt/keybase /usr/share/icons/hicolor
