# etc/systemd/system/fix-resolved.service
[Unit]
Description=Fix /etc/resolv.conf generated by systemd-resolved
After=systemd-resolved.service network.target

[Service]
Type=oneshot
ExecStart=/sbin/restorecon -v /var/run/systemd/resolve/resolv.conf
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
MemoryDenyWriteExecute=yes
ReadWritePaths=/var/run/systemd/resolve/resolv.conf

[Install]
WantedBy=multi-user.target
